# luci-app/etc/init.d/autoedu
#!/bin/sh /etc/rc.common
# Copyright (C) 2025 Matsumiko

START=95
STOP=10

USE_PROCD=1

SCRIPT_DIR="/usr/share/autoedu"
SCRIPT_FILE="$SCRIPT_DIR/auto_edu.py"
ENV_FILE="/root/Auto-Edu/auto_edu.env"
LOG_FILE="/tmp/auto_edu.log"

start_service() {
	local enabled mode interval
	
	config_load autoedu
	config_get enabled config enabled 0
	config_get mode config mode 'EFFICIENT'
	
	# Check if enabled
	[ "$enabled" -eq 0 ] && {
		logger -t autoedu "Service disabled in config, stopping..."
		# Make sure no cron is running
		crontab -l 2>/dev/null | grep -v "auto_edu.py" | grep -v "AUTO_EDU_ENV" | crontab - 2>/dev/null || true
		/etc/init.d/cron restart >/dev/null 2>&1
		return 1
	}
	
	# Check if script exists
	[ ! -f "$SCRIPT_FILE" ] && {
		logger -t autoedu "Script not found: $SCRIPT_FILE"
		return 1
	}
	
	# Sync UCI to .env
	/usr/share/autoedu/sync_config.sh
	
	# Determine cron interval based on mode
	if [ "$mode" = "AGGRESSIVE" ]; then
		interval="*/1 * * * *"
	else
		interval="*/3 * * * *"
	fi
	
	# Setup cron job with duplicate prevention
	logger -t autoedu "Setting up cron: $interval (Mode: $mode)"
	
	# CRITICAL: Remove ALL auto_edu cron entries first (prevent duplicates)
	local temp_cron=$(mktemp)
	crontab -l 2>/dev/null | grep -v "auto_edu.py" | grep -v "AUTO_EDU_ENV" > "$temp_cron" || true
	
	# Add single new entry with lock
	echo "$interval [ ! -f /tmp/auto_edu_python.lock ] && AUTO_EDU_ENV=$ENV_FILE /usr/bin/python3 $SCRIPT_FILE >> $LOG_FILE 2>&1" >> "$temp_cron"
	
	# Apply cron atomically
	crontab "$temp_cron"
	rm -f "$temp_cron"
	
	# Restart cron
	/etc/init.d/cron restart >/dev/null 2>&1
	
	# Verify single entry
	local cron_count=$(crontab -l 2>/dev/null | grep -c "auto_edu.py" || echo 0)
	if [ "$cron_count" -ne 1 ]; then
		logger -t autoedu "WARNING: Expected 1 cron entry, found $cron_count"
	fi
	
	logger -t autoedu "Service started (Mode: $mode, Interval: $interval)"
}

stop_service() {
	logger -t autoedu "Stopping service"
	
	# Remove ALL auto_edu cron entries (even duplicates)
	local temp_cron=$(mktemp)
	crontab -l 2>/dev/null | grep -v "auto_edu.py" | grep -v "AUTO_EDU_ENV" > "$temp_cron" || true
	crontab "$temp_cron"
	rm -f "$temp_cron"
	
	# Restart cron to apply changes
	/etc/init.d/cron restart >/dev/null 2>&1
	
	# Verify removal
	local cron_count=$(crontab -l 2>/dev/null | grep -c "auto_edu.py" || echo 0)
	if [ "$cron_count" -eq 0 ]; then
		logger -t autoedu "All cron entries removed successfully"
	else
		logger -t autoedu "WARNING: Still found $cron_count cron entries after stop"
	fi
	
	# Stop API server
	killall -9 python3 2>/dev/null || true
	
	logger -t autoedu "Service stopped"
}

service_triggers() {
	procd_add_reload_trigger "autoedu"
}

reload_service() {
	logger -t autoedu "Reloading service"
	stop
	start
}
